# salt/roots/salt/docs/supervisor.conf
[program:docs]
command=/envs/docs/bin/httpdocs.py --port 5001
stopsignal=INT
user=vagrant
directory=/vagrant
autostart=false
autorestart=false
